# install/com.machealthkit.watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.machealthkit.watch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>INSTALL_PATH_PLACEHOLDER/scripts/mac_watch.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>600</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/dev/null</string>

  <key>StandardErrorPath</key>
  <string>/tmp/com.machealthkit.watch.err</string>
</dict>
</plist>
